# deploy/systemd/lien-automation-schedule.timer
[Unit]
Description=Legacy single-endpoint scheduler example (production uses per-site weekday Mountain-time triggers)

[Timer]
OnCalendar=Mon..Fri *-*-* 07:00:00 America/Denver
OnCalendar=Mon..Fri *-*-* 11:00:00 America/Denver
OnCalendar=Mon..Fri *-*-* 19:00:00 America/Denver
Persistent=true
Unit=lien-automation-schedule.service

[Install]
WantedBy=timers.target
