# deploy/systemd/lien-automation-schedule.timer
[Unit]
Description=Two daily scheduled scrape triggers (07:30 and 19:30 America/New_York)

[Timer]
OnCalendar=*-*-* 07:30:00 America/New_York
OnCalendar=*-*-* 19:30:00 America/New_York
Persistent=true
Unit=lien-automation-schedule.service

[Install]
WantedBy=timers.target
